# vault/letsencrypt-plugin/docker/crontab
9 4,16 * * * root certbot renew >/proc/1/fd/1 2>/proc/1/fd/2
