# openstack-tenant/packages/mobiledgex/mobiledgex.service
[Unit]
Description=Mobiledgex init service

[Service]
ExecStart=/usr/local/bin/mobiledgex-init.sh

[Install]
WantedBy=default.target
