# openstack-tenant/packages/mobiledgex/mobiledgex.service
[Unit]
Description=Mobiledgex init service
After=open-vm-tools.service
[Service]
ExecStart=/usr/local/bin/mobiledgex-init.sh

[Install]
WantedBy=default.target
