# openstack-tenant/packages/mobiledgex/open-vm-tools.service
[Unit]
Description=Service for virtual machines hosted on VMware
Documentation=http://open-vm-tools.sourceforge.net/about.php
ConditionVirtualization=vmware
DefaultDependencies=no
Requires=dbus.socket
After=dbus.socket
[Service]
ExecStart=/usr/bin/vmtoolsd
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
